# point/setup/etc/supervisor.d/00-prosody.conf
[program:prosody]
command = /usr/bin/lua5.1 /usr/bin/prosody
user = prosody
autostart = true
