# point/setup/etc/supervisor.d/00-redis.conf
[program:redis]
command = /usr/bin/redis-server
user = root
autostart = true
stopsignal = QUIT
