# point/setup/etc/supervisor.d/point-imgproc.conf
[program:point-imgproc]
command = /home/point/venv/bin/python run.py
directory = /home/point/imgproc
user = www-data
autostart = true
