# point/setup/etc/supervisor.d/point-www.conf
[program:point-www]
command = /home/point/venv/bin/geweb run
directory = /home/point/www
user = www-data
autostart = true
