# point/setup/etc/supervisor.d/point-xmpp.conf
[program:point-xmpp]
command = /home/point/venv/bin/python run.py
directory = /home/point/xmpp
user = www-data
autostart = true
