# point/setup/etc/supervisor.d/point-websocket.conf
[program:point-websocket]
command = /home/point/venv/bin/python websocket.py
directory = /home/point/www
user = www-data
autostart = true
